# nvd_ingestion_service/crontab
0 */4 * * * /usr/local/bin/python /app/ingest_nvd_data.py > /proc/1/fd/1 2>/proc/1/fd/2
